# root/etc/cron.daily/nethserver-inventory
#!/bin/bash

exec timeout 120 /usr/sbin/nethserver-inventory
