# root/etc/cron.daily/nethserver-inventory
#!/bin/bash

exec /usr/sbin/nethserver-inventory
